# upstart/asterank.conf
description "Asterank"
start on filesystem and net-device-up IFACE=eth0
stop on shutdown

exec sudo -u ian "/home/ian/asterank/upstart/upstart.sh 2>>/var/log/upstart/asterank.log"
